# uploadfile.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2024-04-13T05:14:52.5860743</Date>
    <Author>DESKTOP-\Dosec101</Author>
    <Description>upload users file to s3 bucket</Description>
    <URI>\upload file</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>Sxxxxxxxxxxxxxxxxxx6-xxxxxxxxxxxxxx-xxxxxxxxxxxxx</UserId>
      <LogonType>InteractiveToken</LogonType>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <LogonTrigger />
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>aws</Command>
      <Arguments>s3 cp C:\Users\Dosec101\Desktop\users.txt s3://testhamadatest1</Arguments>
    </Exec>
  </Actions>
</Task>
